# autostart/clipit-startup.desktop
[Desktop Entry]
Name=ClipIt
Comment=Clipboard Manager
Icon=clipit-trayicon
Exec=clipit
Terminal=false
Type=Application
OnlyShowIn=i3;
