# autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
Comment=Customizable and lightweight notification-daemon
Exec=dunst
Terminal=false
Type=Application
StartupNotify=false
