# zsh/profile.d/beryllium
autoload -Uz add-zsh-hook

PATH="$PATH:$HOME/.local/opt/flutter/bin"
export ANDROID_HOME="$HOME"/.local/share/android-sdk
