# systemd/user/backup-zhistfile.service
[Unit]
Description=backup my zsh history file

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-zhistfile
